# scripts/local.vnc-resolution-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.vnc-resolution-watcher</string>

    <key>ProgramArguments</key>
    <array>
      <string>__SCRIPT_PATH__</string>
      <string>--config</string>
      <string>__CONFIG_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/vnc-resolution-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/vnc-resolution-watcher.err</string>
  </dict>
</plist>
